# scripts/vehicle-daemon.service
[Unit]
Description=Hermes Vehicle Integration Daemon
Documentation=https://github.com/drwjkirkpatrick-web/embedded-vehicle-integration
After=network.target sound.target

[Service]
Type=notify
# Working directory
WorkingDirectory=/home/pi/embedded-vehicle-integration

# Environment
Environment="VEHICLE_LOG_LEVEL=INFO"
Environment="OBD_PORT=/dev/ttyUSB0"
Environment="GPS_PORT=/dev/ttyACM0"
Environment="TG_BOT_TOKEN=REPLACE_ME"
# Load additional env from file (create this with secrets)
EnvironmentFile=-/etc/embedded-vehicle/environment

# Run as pi user (not root for GPIO safety)
User=pi
Group=pi

# Python path — adjust if using venv
ExecStart=/home/pi/.local/bin/uv run --python 3.11 vehicle-daemon
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID

# Graceful shutdown
TimeoutStopSec=90
KillSignal=SIGTERM

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536

# Security hardening (optional — remove if causing issues with GPIO)
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/mnt/nvme0n1p2 /var/log/embedded-vehicle

[Install]
WantedBy=multi-user.target
